# codev/aws/rest-guard/rest-guard.service
[Unit]
Description=Rest Guard

[Service]
Type=oneshot
ExecStart=/{{ORG_NAME}}/rest-guard/rest.sh
